# airootfs/etc/systemd/system/setup.target
[Unit]
Description=Setup Environment
Requires=basic.target
After=basic.target

[Install]
Alias=default.target